# etc/services.d/asterisk/run
#!/usr/bin/execlineb -P
s6-setuidgid asterisk
/usr/sbin/asterisk -f
