# etc/services.d/mariadb/run
#!/usr/bin/execlineb -P
mysqld --user=mysql
