# etc/services.d/mariadb/run
#!/usr/bin/sh
test -e /var/run/mysqld || install -m 755 -o mysql -g root -d /var/run/mysqld
mysqld --user=mysql
