# supervisord.conf
[supervisord]
nodaemon=true

[program:named]
command=/usr/sbin/named -f
